# package/debian/DEBIAN/prerm
#!/bin/bash
systemctl stop tpflow.service || true
systemctl disable tpflow.service || true
