# package/debian/DEBIAN/prerm
#!/bin/bash
systemctl stop tpflow.service
systemctl disable tpflow.service
